# cosmic/live-overlay/etc/systemd/greetd.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin veloxos --noclear %I $TERM
Type=simple
